# services/pickcel-ip.service
[Unit]
Description=fetch local ip address
After=network-online.target

[Service]
Type=oneshot
ExecStart=/usr/bin/bash /.pickcel/scripts/ip.sh

[Install]
WantedBy=multi-user.target